# deploy/systemd/fourok-postgres-backup.service
[Unit]
Description=4ok PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/4ok
Environment=GCB_IMAGE_TAG=replace-with-commit-hash
Environment=GCB_CONFIG_PATH=/opt/4ok/.local/4ok.toml
EnvironmentFile=/etc/4ok/4ok.env
ExecStart=/bin/sh -lc '/usr/bin/docker compose run --rm app postgres-backup --database-url "$GCB_DATABASE_URL" --output "/var/lib/4ok/backups/4ok-$(date +%%Y%%m%%d-%%H%%M%%S).dump"'
