# deploy/systemd/fourok-retention.service
[Unit]
Description=4ok retention purge
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/4ok
Environment=GCB_IMAGE_TAG=replace-with-commit-hash
Environment=GCB_CONFIG_PATH=/opt/4ok/.local/4ok.toml
EnvironmentFile=/etc/4ok/4ok.env
ExecStart=/usr/bin/docker compose run --rm app purge-raw-retention --database-url ${GCB_DATABASE_URL} --config /etc/4ok/4ok.toml
ExecStart=/usr/bin/docker compose run --rm app purge-audit-retention --database-url ${GCB_DATABASE_URL} --config /etc/4ok/4ok.toml
ExecStart=/usr/bin/docker compose run --rm app purge-webhook-retention --database-url ${GCB_DATABASE_URL} --config /etc/4ok/4ok.toml
ExecStart=/usr/bin/docker compose run --rm app purge-backup-retention --config /etc/4ok/4ok.toml
