# deploy/systemd/fourok-retention.timer
[Unit]
Description=Run fourok retention purge daily

[Timer]
OnCalendar=*-*-* 03:15:00
AccuracySec=10min
Persistent=true
Unit=fourok-retention.service

[Install]
WantedBy=timers.target
