# deploy/systemd/fourok-postgres-backup.timer
[Unit]
Description=Run 4ok PostgreSQL backup nightly

[Timer]
OnCalendar=*-*-* 02:15:00
AccuracySec=5min
Persistent=true
Unit=4ok-postgres-backup.service

[Install]
WantedBy=timers.target
