# deploy/systemd/fourok-retention.timer
[Unit]
Description=Run 4ok retention purge daily

[Timer]
OnCalendar=*-*-* 03:15:00
AccuracySec=10min
Persistent=true
Unit=4ok-retention.service

[Install]
WantedBy=timers.target
