# deploy/systemd/fourok-postgres-backup.timer
[Unit]
Description=Run fourok PostgreSQL backup nightly

[Timer]
OnCalendar=*-*-* 02:15:00
AccuracySec=5min
Persistent=true
Unit=fourok-postgres-backup.service

[Install]
WantedBy=timers.target
